# supervisord.conf
; supervisor config file

[unix_http_server]
file=/var/run/supervisor.sock   ; (the path to the socket file)
chmod=0700                       ; sockef file mode (default 0700)

[supervisord]
logfile=/var/log/supervisor/supervisord.log ; (main log file;default $CWD/supervisord.log)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
childlogdir=/var/log/supervisor            ; ('AUTO' child log dir, default $TEMP)

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket

; The [include] section can just contain the "files" setting.  This
; setting can list multiple files (separated by whitespace or
; newlines).  It can also contain wildcards.  The filenames are
; interpreted as relative to this file.  Included files *cannot*
; include files themselves.

[include]
files = /etc/supervisor/conf.d/*.conf

[program:backend-sad-jurangjero]
command=/home/informa/.local/share/virtualenvs/jurangjero-sB0ecPrY/bin/gunicorn --bind 0.0.0.0:8080 api_sad_sig.wsgi
directory=/var/www/backend-sad-jurangjero
user=informa
autostart=true
autorestart=true
redirect_stderr=true

[program:backend-sad-blanciran]
command=/home/informa/.local/share/virtualenvs/blanciran-KGrf4nbc/bin/gunicorn --bind 0.0.0.0:8081 api_sad_sig.wsgi
directory=/home/informa/be-app/blanciran
user=informa
autostart=true
autorestart=true
redirect_stderr=true

[program:backend-sad-brangkal]
command=/home/informa/.local/share/virtualenvs/brangkal-CT6FihBz/bin/gunicorn --bind 0.0.0.0:8082 api_sad_sig.wsgi
directory=/home/informa/be-app/brangkal
user=informa
autostart=true
autorestart=true
redirect_stderr=true

[program:backend-sad-gempol]
command=/home/informa/.local/share/virtualenvs/gempol-9c7CrGXy/bin/gunicorn --bind 0.0.0.0:8083 api_sad_sig.wsgi
directory=/home/informa/be-app/gempol
user=informa
autostart=true
autorestart=true
redirect_stderr=true

[program:backend-sad-gledeg]
command=/home/informa/.local/share/virtualenvs/gledeg-c608xCtL/bin/gunicorn --bind 0.0.0.0:8084 api_sad_sig.wsgi
directory=/home/informa/be-app/gledeg
user=informa
autostart=true
autorestart=true
redirect_stderr=true

[program:backend-sad-jambeyan]
command=/home/informa/.local/share/virtualenvs/jambeyan-wHPN6DXZ/bin/gunicorn --bind 0.0.0.0:8085 api_sad_sig.wsgi
directory=/home/informa/be-app/jambeyan
user=informa
autostart=true
autorestart=true
redirect_stderr=true

[program:backend-sad-kadirejo]
command=/home/informa/.local/share/virtualenvs/kadirejo-8ctZEJYR/bin/gunicorn --bind 0.0.0.0:8086 api_sad_sig.wsgi
directory=/home/informa/be-app/kadirejo
user=informa
autostart=true
autorestart=true
redirect_stderr=true

[program:backend-sad-karangan]
command=/home/informa/.local/share/virtualenvs/karangan-hsiyeduN/bin/gunicorn --bind 0.0.0.0:8087 api_sad_sig.wsgi
directory=/home/informa/be-app/karangan
user=informa
autostart=true
autorestart=true
redirect_stderr=true

[program:backend-sad-karanganom]
command=/home/informa/.local/share/virtualenvs/karanganom-bhZlOE_M/bin/gunicorn --bind 0.0.0.0:8088 api_sad_sig.wsgi
directory=/home/informa/be-app/karanganom
user=informa
autostart=true
autorestart=true
redirect_stderr=true

[program:backend-sad-kunden]
command=/home/informa/.local/share/virtualenvs/kunden-4I1hmz3w/bin/gunicorn --bind 0.0.0.0:8089 api_sad_sig.wsgi
directory=/home/informa/be-app/kunden
user=informa
autostart=true
autorestart=true
redirect_stderr=true

[program:backend-sad-ngabeyan]
command=/home/informa/.local/share/virtualenvs/ngabeyan-xiGUm0Vv/bin/gunicorn --bind 0.0.0.0:8090 api_sad_sig.wsgi
directory=/home/informa/be-app/ngabeyan
user=informa
autostart=true
autorestart=true
redirect_stderr=true

[program:backend-sad-padas]
command=/home/informa/.local/share/virtualenvs/padas-_rFqutg1/bin/gunicorn --bind 0.0.0.0:8091 api_sad_sig.wsgi
directory=/home/informa/be-app/padas
user=informa
autostart=true
autorestart=true
redirect_stderr=true

[program:backend-sad-pondok]
command=/home/informa/.local/share/virtualenvs/pondok-ZzrMgSPg/bin/gunicorn --bind 0.0.0.0:8092 api_sad_sig.wsgi
directory=/home/informa/be-app/pondok
user=informa
autostart=true
autorestart=true
redirect_stderr=true

[program:backend-sad-tarubasan]
command=/home/informa/.local/share/virtualenvs/tarubasan-2R0MSFfF/bin/gunicorn --bind 0.0.0.0:8093 api_sad_sig.wsgi
directory=/home/informa/be-app/tarubasan
user=informa
autostart=true
autorestart=true
redirect_stderr=true

[program:backend-sad-troso]
command=/home/informa/.local/share/virtualenvs/troso-Qd6_BIEA/bin/gunicorn --bind 0.0.0.0:8094 api_sad_sig.wsgi
directory=/home/informa/be-app/troso
user=informa
autostart=true
autorestart=true
redirect_stderr=true

[program:backend-sad-soropaten]
command=/home/informa/.local/share/virtualenvs/soropaten-PLMUylXX/bin/gunicorn --bind 0.0.0.0:8095 api_sad_sig.wsgi
directory=/home/informa/be-app/soropaten
user=informa
autostart=true
autorestart=true
redirect_stderr=true

[program:backend-sad-jeblok]
command=/home/informa/.local/share/virtualenvs/jeblok-eo73BAkU/bin/gunicorn --bind 0.0.0.0:8096 api_sad_sig.wsgi
directory=/home/informa/be-app/jeblok
user=informa
autostart=true
autorestart=true
redirect_stderr=true

[program:backend-sad-jungkare]
command=/home/informa/.local/share/virtualenvs/jungkare-o4N0WD3s/bin/gunicorn --bind 0.0.0.0:8097 api_sad_sig.wsgi
directory=/home/informa/be-app/jungkare
user=informa
autostart=true
autorestart=true
redirect_stderr=true

[program:backend-sad-beku]
command=/home/informa/.local/share/virtualenvs/beku-NUW9_pFh/bin/gunicorn --bind 0.0.0.0:8098 api_sad_sig.wsgi
directory=/home/informa/be-app/beku
user=informa
autostart=true
autorestart=true
redirect_stderr=true

[program:backend-sad-demo]
command=/home/informa/.local/share/virtualenvs/demo-9XkVDRW1/bin/gunicorn --bind 0.0.0.0:8099 api_sad_sig.wsgi
directory=/home/informa/be-app/demo
user=informa
autostart=true
autorestart=true
redirect_stderr=true